# deploy/supervisord/supervisord-redis.conf
[program:redis]
command=/usr/local/bin/redis-server /plex-utils/deploy/redis.conf
directory=/redis/redis_data
user=root
autostart=true
stopsignal=QUIT