# deploy/supervisord/supervisord-plex-utils.conf
[program:celery_plex-utils]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=0
numprocs=1
user=www-data
directory=/plex-utils
priority=300
startinorder=true
command=/plex-utils/bin/celery -A plex_utils worker -l info -Q plex-utils -c 4 -E -n plex-utils@%%h
EnvironmentFile=/config/.env