# deploy/supervisord/supervisord-uwsgi.conf
[program:uwsgi]
command=/plex-utils/bin/uwsgi --ini /plex-utils/deploy/uwsgi.ini
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=100
startinorder=true
startsecs=0
EnvironmentFile=/config/.env