# app/static/dockerfiles/plex-utills.conf
[program:main]
command=/app/venv/bin/gunicorn -b localhost:5000 -w 4 main:app
directory=/app/plex-utills
user=plex-utils
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true